# rootfs/service/rtorrent/finish
#!/bin/sh

echo
echo "rTorrent exited."

# Send SIGTERM to all remaining services
s6-svscanctl -t /service
