# rootfs/service/rtorrent/finish
#!/bin/sh

echo
echo "rTorrent exited."
s6-svscanctl -t /service
sleep 3
