# rootfs/service/.s6-svscan/finish
#!/bin/bash

echo
echo "Shutting down..."

# Kind of workaround...
touch /etc/leapsecs.dat
exec s6-svwait -d /service/rtorrent
