# rootfs/etc/s6/.s6-svscan/finish
#!/bin/sh

echo "[docktorrent] Shutting down..."
