# rootfs/etc/s6/.s6-svscan/finish
#!/bin/bash

echo
echo "Shutting down..."

# There must have a better way than sleep
sleep 3
