# rootfs/etc/s6/rtorrent/finish
#!/bin/sh

echo
echo "rTorrent exited."
s6-svscanctl -t /etc/s6
sleep 3
